# etc/cron.daily/ar-sync
#!/bin/bash

/usr/libexec/ar-sync/poem-sync
/usr/libexec/ar-sync/topology-sync

